# crontab
00 12 * * 1-5        python2 /home/tux/Projets/timesheet/auto-fill-openerp-obj.py
